# cron/crontab
# This file created automatically
# Do not edit directly (your changes will be overwritten)
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/test/exampleCron.py
