# cron/crontab
# This file created automatically
# Do not edit directly (your changes will be overwritten)
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testhandler/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 rooat /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/testQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thequiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thequiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/theQuiz2/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/boringquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/myquiz1/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/thefunquiz/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bloodmagic/exampleCron.py
# contents of this file auto-appended to cron/crontab during build
* 0 0 0 0 root /scripts/bonegerbil/exampleCron.py
